# update-from-added-sources.path
[Unit]
Description=Watch sources.list.d for added package sources

[Path]
PathChanged=/etc/apt/sources.list.d

[Install]
WantedBy=multi-user.target
